# crontab
@reboot /usr/local/bin/start_server
*/5 * * * * /usr/local/bin/save_server
